# root/etc/init.d/minigate
#!/bin/sh /etc/rc.common
START=99
STOP=10

start() {
    local en=$(uci -q get minigate.global.enabled)
    [ "$en" = "1" ] || { logger -t minigate "未启用"; return 0; }
    logger -t minigate "启动..."
    /bin/sh /usr/lib/minigate/proxy.sh start
    /bin/sh /usr/lib/minigate/ddns.sh force &
    # 设置 cron
    sed -i '/minigate/d' /etc/crontabs/root 2>/dev/null
    echo "*/5 * * * * /bin/sh /usr/lib/minigate/ddns.sh" >> /etc/crontabs/root
    echo "0 3 * * * /bin/sh /usr/lib/minigate/acme.sh renew" >> /etc/crontabs/root
    echo "*/2 * * * * /bin/sh /usr/lib/minigate/geofence.sh scan" >> /etc/crontabs/root
    /etc/init.d/cron restart 2>/dev/null
    logger -t minigate "已启动"
}

stop() {
    logger -t minigate "停止所有服务..."
    # 停止反向代理
    /bin/sh /usr/lib/minigate/proxy.sh stop
    sleep 1
    # 移除所有 cron 定时任务（DDNS + ACME）
    sed -i '/minigate/d' /etc/crontabs/root 2>/dev/null
    /etc/init.d/cron restart 2>/dev/null
    # 清除 DDNS 运行状态
    local sections=$(uci -q show minigate | grep '=ddns$' | cut -d. -f2 | cut -d= -f1)
    for sec in $sections; do
        uci -q set minigate.${sec}.status="stopped"
        uci -q set minigate.${sec}.status_msg="已停止"
        uci -q set minigate.${sec}.next_sync=""
    done
    uci commit minigate
    logger -t minigate "已停止（DDNS定时/ACME续期/反代 全部关闭）"
}

restart() {
    stop
    sleep 1
    start
}

reload() {
    /bin/sh /usr/lib/minigate/proxy.sh reload
    /bin/sh /usr/lib/minigate/ddns.sh force
}
